# systemd/gb-mcp-servers.timer
[Unit]
Description=Daily MCP server update check

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
